# mac-service/com.workshop.labelservice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent to auto-start the Workshop Label Service at login and keep it
  running. Install with install-autostart.sh (which rewrites the paths below
  for your machine).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.workshop.labelservice</string>
  <key>ProgramArguments</key>
  <array>
    <string>__RUN_SH__</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/launchd.err.log</string>
</dict>
</plist>
